# deploy/corerp.service
[Unit]
Description=CoreRP Persistent Narrative Runtime
After=network.target

[Service]
Type=simple
User=kelebit
WorkingDirectory=/home/kelebit/corerp
ExecStart=/home/kelebit/corerp/corerp serve \
  -characters /home/kelebit/corerp/characters \
  -llm-url $LLM_URL \
  -llm-key $LLM_API_KEY \
  -llm-model $LLM_MODEL
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=LLM_URL=https://gcli.ggchan.dev/v1
Environment=LLM_MODEL=gemini-3-flash-preview
Environment=LLM_API_KEY=

# Fill in LLM_API_KEY above or use a separate env file:
# EnvironmentFile=/home/kelebit/corerp/.env

# SQLite backup hint: cp data/memory.db backup/memory_$(date +%Y%m%d).db

[Install]
WantedBy=multi-user.target
